# supplychains-bot.service
[Unit]
Description=SupplyChains Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/supplychains-bot
ExecStart=/usr/bin/python3 /opt/supplychains-bot/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
